# templates/dispatch-api.service.tmpl
[Unit]
Description=Agent GTD Dispatch API
After=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{WORKING_DIR}}
EnvironmentFile={{ENV_FILE}}
ExecStart=/home/{{SERVICE_USER}}/.local/bin/agent-gtd-dispatch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
